# cyberplasma/systemd/user/bismuth-mode.service
# Restore Bismuth tiling mode at login.
# Runs the apply_bismuth_mode script if a saved mode exists.

[Unit]
Description=Restore Bismuth tiling mode
After=graphical-session.target

[Service]
ExecStart=/usr/bin/env apply_bismuth_mode.sh

[Install]
WantedBy=cyberplasma.target
