# cyberplasma/systemd/user/eww.service
# CyberPlasma EWW widgets service.
# Requires the `eww` executable to be installed.
# Limit permissions and avoid exposing sensitive paths or data.

[Unit]
Description=CyberPlasma EWW widgets
After=graphical-session.target

[Service]
ExecStart=eww daemon
ExecStartPost=eww open cyberbar
Restart=on-failure

[Install]
WantedBy=cyberplasma.target
